# src/piper_teleop/web_server/nginx/tactile-teleop.service
[Unit]
Description=Tactile Teleop FastAPI Web Server
After=network.target
Wants=network.target

[Service]
Type=exec
User=zhamers
Group=zhamers
WorkingDirectory=/home/zhamers/telegrip

# Environment setup - activate conda environment
Environment=PATH=/home/zhamers/miniconda3/envs/tactile-teleop/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=CONDA_DEFAULT_ENV=tactile-teleop
Environment=CONDA_PREFIX=/home/zhamers/miniconda3/envs/tactile-teleop

# Start the FastAPI server in proxy mode with Unix socket
ExecStart=/home/zhamers/miniconda3/envs/tactile-teleop/bin/python -m piper_teleop.web_server --behind-proxy --log-level info

# Restart policy
Restart=always
RestartSec=5
StartLimitBurst=3
StartLimitInterval=60

# Process management
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30

# Security (optional hardening)
NoNewPrivileges=yes
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tactile-teleop

[Install]
WantedBy=multi-user.target